# roles/utils/templates/etc/profile.d/history.sh.j2
# {{ ansible_managed }}
shopt -s histappend
HISTSIZE=100000
HISTFILESIZE=200000
HISTTIMEFORMAT="%F - %T : "
HISTCONTROL=ignoredups
PROMPT_COMMAND="history -a;history -c;history -r;$PROMPT_COMMAND"
