# ansible/roles/utils/templates/etc/cron.daily/logrotate.j2
# {{ ansible_managed }}

#!/usr/bin/env bash
exit 0
