# roles/utils/templates/etc/cron.daily/logrotate.j2
# {{ ansible_managed }}

#!/bin/sh
exit 0
