# ansible/roles/utils/templates/etc/profile.d/alias.sh.j2
# {{ ansible_managed }}

alias ls='ls --color=auto'
